# lib/capistrano/templates/launchd/deliver.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>de.faxomat.deliver</string>
  <key>ProgramArguments</key>
  <array>
    <string>bundle</string>
    <string>exec</string>
    <string>rake</string>
    <string>faxomat:deliver</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/local/opt/ruby/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>RAILS_ENV</key>
    <string>production</string>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StartCalendarInterval</key>
  <array>
    <% (1..5).each do |weekday| %>
      <% (8..17).each do |hour| %>
        <dict>
          <key>Weekday</key>
          <integer><%= weekday %></integer>
          <key>Hour</key>
          <integer><%= hour %></integer>
          <key>Minute</key>
          <integer>15</integer>
        </dict>
      <% end %>
    <% end %>
  </array>
  <key>WorkingDirectory</key>
  <string><%= current_path %></string>
  <key>UserName</key>
  <string><%= fetch(:user) %></string>
  <key>GroupName</key>
  <string>staff</string>
  <key>StandardOutPath</key>
  <string><%= current_path.join('log', 'deliver.log') %></string>
  <key>StandardErrorPath</key>
  <string><%= current_path.join('log', 'deliver.log') %></string>
</dict>
</plist>
